# newsbot.service
[Unit]
Description=AI Newsletter Bot
After=network.target

[Service]
Type=simple
User=newsbot
Group=newsbot
WorkingDirectory=/opt/ai_newsletter
Environment=DISPLAY=:99
ExecStartPre=/usr/bin/Xvfb :99 -screen 0 1024x768x24 -ac
ExecStart=/opt/ai_newsletter/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/newsbot/output.log
StandardError=append:/var/log/newsbot/error.log

[Install]
WantedBy=multi-user.target 